# config/discord-indexer-bot.service
[Unit]
Description=Discord Indexer Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/discord-indexer
Environment=PYTHONPATH=/root/discord-indexer
ExecStart=/root/discord-indexer/venv/bin/python /root/discord-indexer/src/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-indexer-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/root/discord-indexer

[Install]
WantedBy=multi-user.target